# system_units/currentyear.timer
# /etc/systemd/system/currentyear.timer

[Unit]
Description=The current year, hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target

